# deploy/systemd/signalmaker-pipeline.service
[Unit]
Description=SignalMaker Pipeline Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/signalmaker
ExecStart=/bin/bash /opt/signalmaker/scripts/start_pipeline_worker.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
