# deploy/systemd/signalmaker-frontend.service
[Unit]
Description=SignalMaker Frontend UI
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/pi/Desktop/SignalMaker
ExecStart=/bin/bash /home/pi/Desktop/SignalMaker/scripts/start_frontend.sh
Restart=on-failure
RestartSec=5
Environment=FRONTEND_PORT=8090
Environment=SIGNALMAKER_API_BASE=http://127.0.0.1:5000

[Install]
WantedBy=multi-user.target
