# deploy/systemd/signalmaker-api.service
[Unit]
Description=SignalMaker API
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/signalmaker
ExecStart=/bin/bash /opt/signalmaker/scripts/start_api.sh
Restart=always
RestartSec=5
Environment=APP_PORT=5000

[Install]
WantedBy=multi-user.target
