# systemd/raspberry-executor.service
[Unit]
Description=SignalMaker Raspberry Executor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Desktop/SignalMaker
Environment=PYTHONUNBUFFERED=1
ExecStart=/bin/bash /home/pi/Desktop/SignalMaker/run_bot_service.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
