# systemd/signalmaker-kiosk.service
[Unit]
Description=SignalMaker Raspberry web kiosk
After=network-online.target graphical.target signalmaker-api.service
Wants=network-online.target graphical.target signalmaker-api.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Desktop/SignalMaker
Environment=APP_PORT=8080
Environment=SIGNALMAKER_KIOSK_URL=http://127.0.0.1:8080/index.html
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
ExecStart=/bin/bash -lc 'exec /home/pi/Desktop/SignalMaker/kiosk.sh'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target
